# ops/systemd/resource-provider-gpu.service
[Unit]
Description=CoopHive V0 Resource Provider GPU
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]
Environment="LOG_TYPE=json"
Environment="DEBUG=true"
Environment="HOME=/app/coophive"
Environment="OFFER_GPU=1"
EnvironmentFile=/app/coophive/.env
Restart=always
RestartSec=5s
ExecStart=/usr/bin/hive resource-provider

[Install]
WantedBy=multi-user.target
